# scheduling/user_level/capslock-fix.timer
[Unit]
Description=Run capslock-fix.service 15 seconds after the user logs in

[Timer]
OnActiveSec=10

[Install]
WantedBy=timers.target
