# scheduling/user_level/capslock-fix.timer
[Unit]
Description=Run capslock-fix.service 30 seconds after the user logs in

[Timer]
OnStartupSec=30

[Install]
WantedBy=timers.target

