# scheduling/system_level/distro_size.service
[Unit]
Description=Check distro size once a day

[Service]
Type=oneshot
ExecStart=/home/ali/system_scripts/disc_usage/distro_size.sh --systemd
