# scheduling/system_level/distro_size.timer
[Unit]
Description=Run distro_size.service 2 minutes after the user logs in

[Timer]
OnStartupSec=2min

[Install]
WantedBy=timers.target
